# scripts/com.daglo.youtube-pipeline.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.daglo.youtube-pipeline</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/bez/development/daglo-mcp/scripts/youtube-auto-pipeline.mjs</string>
      <string>/Users/bez/development/daglo-mcp/scripts/youtube-pipeline.config.json</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>DAGLO_EMAIL</key>
      <string>your-email@example.com</string>
      <key>DAGLO_PASSWORD</key>
      <string>your-password</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/bez/development/daglo-mcp</string>

    <key>StandardOutPath</key>
    <string>/Users/bez/development/daglo-mcp/.work/youtube-pipeline/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bez/development/daglo-mcp/.work/youtube-pipeline/launchd.err.log</string>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
